# 90-ambient-backlight-control.rules
# udev rules for ambient-backlight-control
# Grants write access to display and keyboard backlight sysfs nodes
# without requiring root at runtime.
#
# Install: sudo cp 90-ambient-backlight-control.rules /etc/udev/rules.d/
# Reload:  sudo udevadm control --reload-rules && sudo udevadm trigger

SUBSYSTEM=="leds",      KERNEL=="platform::kbd_backlight", RUN+="/bin/chmod 0666 /sys%p/brightness"
SUBSYSTEM=="backlight", KERNEL=="intel_backlight",         RUN+="/bin/chmod 0666 /sys%p/brightness"
